# scripts/Library/LaunchAgents/com.qiushi.keep-codex-fast-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.qiushi.keep-codex-fast-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/qiushi/.local/bin/keep-codex-fast-biweekly-report</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/qiushi/Library/Logs/CodexMaintenance/keep-codex-fast-launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/qiushi/Library/Logs/CodexMaintenance/keep-codex-fast-launchd.err.log</string>
</dict>
</plist>
