# scripts/Library/LaunchAgents/com.qiushi.caffeinate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.qiushi.caffeinate</string>

  <!-- -i blocks idle system sleep but NOT lid-close sleep, so the laptop
       still sleeps when I shut the lid, but never falls asleep on its own. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-i</string>
  </array>

  <!-- Start at login and respawn if it ever dies, so I never have to
       remember to run caffeinate myself. -->
  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
